# mnt-emu.automount
[Unit]
Description=Samba automount for emu
After=network.target

[Automount]
Where=/mnt/emu
TimeoutIdleSec=180

[Install]
WantedBy=multi-user.target
